# pterodactyl/files/etc/supervisor.d/php-fpm.ini
[program:php-fpm]
command=/usr/sbin/php-fpm7 --nodaemonize -c PANEL_ROOT/etc/php7 -y PANEL_ROOT/etc/php7/php-fpm.conf
autostart=false
autorestart=true
priority=5
