# pterodactyl/files/etc/supervisor.d/php-fpm.ini
[program:php-fpm]
command=/usr/sbin/php-fpm --nodaemonize -c PANEL_ROOT/etc/php -y PANEL_ROOT/etc/php/php-fpm.conf
autostart=false
autorestart=true
priority=5
